# patchup.service
[Unit]
Description=Patchup daemon

[Service]
Type=oneshot
ExecStart=/usr/local/bin/patchup

[Install]
WantedBy=multi-user.target
